# files/supervisord.conf
[supervisord]
nodaemon=true

[program:cassandra]
command=/home/apache-cassandra-2.0.7/bin/cassandra -f

